# debian/centrum-djangobaselibrary-djangobaseproject-auto.postinst
#!/bin/bash

set -e

BASE=djangobaselibrary
APP=djangobaseproject

APACHE_CONFIG="$BASE-$APP-mod-wsgi"
APP_DIR="/centrum/$BASE/$APP"
CONFIG_DIR="/etc/djangobaselibrary/"

case "$1" in
    configure)

    # byte compiling
    PYCOMPILE=$( python -c 'print __import__("compileall").__file__' || true )
    [[ -f $PYCOMPILE ]] && python $PYCOMPILE $APP_DIR $CONFIG_DIR

    # enable apache modules
    a2enmod proxy_http
    a2enmod deflate
    a2enmod alias

    # enable ego modules
    a2enmod mod_auth_ego
    a2enmod mod_form
    a2enmod mod_uri
    a2enmod mod_cookies
    a2enmod RS3

    # enable mod python site
    a2ensite $APACHE_CONFIG

    # reload apache
    #invoke-rc.d apache2 reload

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


