# debian/centrum-djangobaselibrary-djangobaseproject.cron.cron.d
# cron jobs for centrum-djangobaselibrary-djangobaseproject package
*/5 * * * * root run-parts /centrum/djangobaselibrary/djangobaseproject/cron/every5min
10 2 * * * root run-parts /centrum/djangobaselibrary/djangobaseproject/cron/daily

